# rootfs/etc/s6/aredn-manager/run
#!/bin/sh

exec aredn-manager -d server --no-daemon