# rootfs/etc/s6/api/run
#!/bin/sh

cd /api
exec node src/index.js